# roles/node_monitor/files/vector.service
[Unit]
Description=Vector Log Collector
Documentation=https://vector.dev/docs/
After=network-online.target
Wants=network-online.target

[Service]
User=root
EnvironmentFile=-/etc/default/vector
ExecStart=/usr/bin/vector $VECTOR_OPTS
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5s
LimitNOFILE=655360

[Install]
WantedBy=multi-user.target
